# deploy/violation-worker.service
[Unit]
Description=Violation AI - background classify worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/violation-bot
ExecStart=/root/violation-bot/.venv-webapp/bin/python -m webapp.worker
Restart=always
RestartSec=5
# Default StandardOutput=journal — visible via `journalctl -u violation-worker -f`
# and rotated automatically. See companion unit for context.

[Install]
WantedBy=multi-user.target
